# projects/asset_service/asset_sync.timer
[Unit]
Description=periodically synchronize assets with SOPHIA cloud

[Timer]
OnUnitActiveSec=2h
OnStartupSec=30s

[Install]
WantedBy=timers.target

